# ops/systemd/servora-med.service
[Unit]
Description=Servora-Med API (Fastify)
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=simple
User=servora-med
Group=servora-med
WorkingDirectory=/opt/servora-med/current/server
Environment=NODE_ENV=production
EnvironmentFile=/etc/servora-med/servora-med.env
ExecStart=/usr/bin/node dist/index.js
Restart=on-failure
RestartSec=5
TimeoutStopSec=30
KillMode=mixed
UMask=0077
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true

StandardOutput=journal
StandardError=journal
SyslogIdentifier=servora-med

[Install]
WantedBy=multi-user.target
